# cookbooks/tomee/templates/default/tomee.upstart.conf.erb
#!upstart
description "Tomcat Server"

start on vagrant-mounted MOUNTPOINT=/opt/tomee/lib
stop on runlevel [!2345]
respawn
respawn limit 10 5

# run as non privileged user 
# add user with this command:
## adduser --system --ingroup www-data --home /opt/apache-tomcat apache-tomcat
setuid vagrant
setgid vagrant

# paths:
env CATALINA_BASE="<%= node[:tomee][:dir] %>"
env CATALINA_HOME="<%= node[:tomee][:dir] %>"
env CATALINA_PID="<%= node[:tomee][:dir] %>/catalina.pid"

exec $CATALINA_HOME/bin/catalina.sh start

# cleanup temp directory after stop
post-stop script 
	rm -rf $CATALINA_HOME/temp/*
end script